# etc/maxscale.service.in
[Unit]
Description=MariaDB MaxScale Database Proxy
After=network.target

[Service]
Type=forking
Restart=on-abort

# Make sure /var/run/maxscale exists
PermissionsStartOnly=true
ExecStartPre=/usr/bin/install -d @MAXSCALE_VARDIR@/run/maxscale -o maxscale -g maxscale

PIDFile=@MAXSCALE_VARDIR@/run/maxscale/maxscale.pid

# Use the default user and group
User=maxscale
Group=maxscale

ExecStart=@CMAKE_INSTALL_PREFIX@/@MAXSCALE_BINDIR@/maxscale
TimeoutStartSec=120
LimitNOFILE=65535
StartLimitBurst=0
#Please see the MaxScale release notes and/or watchdog documentation
#before turning the watchdog on.
#WatchdogSec=60s

# Only relevant when MaxScale is linked with -fsanitize=address
Environment=ASAN_OPTIONS=abort_on_error=1

[Install]
WantedBy=multi-user.target
