# etc/maxscale.service.in
[Unit]
Description=MariaDB MaxScale Database Proxy
After=network.target

[Service]
Type=forking
Restart=on-failure
PIDFile=@MAXSCALE_VARDIR@/run/maxscale/maxscale.pid
ExecStartPre=-/usr/bin/mkdir @MAXSCALE_VARDIR@/run/maxscale
ExecStartPre=/usr/bin/chown -R maxscale:maxscale @MAXSCALE_VARDIR@/run/maxscale
ExecStart=@CMAKE_INSTALL_PREFIX@/@MAXSCALE_BINDIR@/maxscale --user=maxscale
LimitNOFILE=4096

[Install]
WantedBy=multi-user.target
