# etc/maxscale.service.in
[Unit]
Description=MariaDB MaxScale Database Proxy
After=network.target

[Service]
Type=forking
Restart=on-abort
@SERVICE_FILE_DEBUG_OPTIONS@
@SERVICE_FILE_ASAN_OPTIONS@
@SERVICE_FILE_UBSAN_OPTIONS@

# This makes sure /var/run/maxscale exists. The directory is created as
# /run/maxscale but /var/run is always a symlink to /run which makes it
# work even with the legacy /var/run that MaxScale uses.
RuntimeDirectory=maxscale
PIDFile=@MAXSCALE_VARDIR@/run/maxscale/maxscale.pid

# Use the default user and group
User=maxscale
Group=maxscale

ExecStart=@CMAKE_INSTALL_PREFIX@/@MAXSCALE_BINDIR@/maxscale
TimeoutStartSec=120
LimitNOFILE=65535
StartLimitBurst=0
WatchdogSec=60s
NotifyAccess=all

# MaxScale should be restarted if it exits with 75 (BSD's EX_TEMPFAIL)
RestartForceExitStatus=75

[Install]
WantedBy=multi-user.target
