# etc/upstart/maxscale.conf.in
# MaxScale service

description "MariaDB MaxScale"

start on stopped rc RUNLEVEL=[2345]
stop on starting rc runlevel [!2345]

# Respawn the process on abnormal exits
respawn

# Uncomment this to limit respawns to two every five minutes
# respawn limit 2 5

# Unlimited open files
limit nofile unlimited unlimited

# Make sure /var/run/maxscale exists
pre-start exec /usr/bin/install -d -o maxscale -g maxscale @MAXSCALE_VARDIR@/run/maxscale

# Change the user to maxscale:maxscale
setuid maxscale
setgid maxscale

# Start MaxScale
expect fork
exec /usr/bin/maxscale
